# libcr4shsend/Package/DEBIAN/prerm
#!/bin/sh

/bin/launchctl unload /Library/LaunchDaemons/com.muirey03.cr4shedd.plist
/bin/launchctl load /Library/LaunchDaemons/com.muirey03.cr4shedd.plist

killall ReportCrash 2> /dev/null


exit 0
